# macos/com.fran.bibliotecario-pf.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fran.bibliotecario-pf</string>

  <!-- Recarga la redirección pf 80 -> 8080 en cada arranque. pf no conserva
       reglas de anclas no-Apple entre reinicios, por eso lo reaplicamos. -->
  <key>ProgramArguments</key>
  <array>
    <string>/sbin/pfctl</string>
    <string>-ef</string>
    <string>/etc/pf-bibliotecario.conf</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/com.fran.bibliotecario-pf.err</string>
</dict>
</plist>
